# deploy/systemd/prox-backup.timer
[Unit]
Description=Create daily encrypted prox production backup

[Timer]
OnCalendar=*-*-* 03:20:00
RandomizedDelaySec=20min
Persistent=true
Unit=prox-backup.service

[Install]
WantedBy=timers.target
