# deploy/systemd/prox-restore-drill.service
[Unit]
Description=Verify the latest prox backup in isolated PostgreSQL
After=docker.service prox-backup.service
Requires=docker.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/prox/current
EnvironmentFile=-/etc/prox/operations.env
ExecStart=/usr/bin/bash /opt/prox/current/scripts/deploy/restore-drill.sh
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadOnlyPaths=/etc/prox
ReadWritePaths=/var/backups/prox /tmp
TimeoutStartSec=30min
